# rootfs/etc/systemd/system/firstboot.service
[Unit]
Description=Take care of simage resizing on first boot

[Service]
Type=oneshot
ExecStart=/firstboot.sh

[Install]
WantedBy=multi-user.target